# pnl-sentinel/deploy/stockpulse/stockpulse-webhook.service
[Unit]
Description=StockPulse webhook (FastAPI/uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=stockpulse
WorkingDirectory=/opt/stockpulse
EnvironmentFile=-/opt/stockpulse/.env
ExecStartPre=/opt/stockpulse/deploy/stockpulse/fetch_secrets.sh
ExecStart=/opt/stockpulse/.venv/bin/python -m uvicorn webhook:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stockpulse-webhook

[Install]
WantedBy=multi-user.target
