# pnl-sentinel/deploy/stockpulse/stockpulse-bot.service
[Unit]
Description=StockPulse Telegram bot (long-poll)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=stockpulse
WorkingDirectory=/opt/stockpulse
EnvironmentFile=/opt/stockpulse/.env
ExecStartPre=/opt/stockpulse/deploy/stockpulse/fetch_secrets.sh
ExecStart=/opt/stockpulse/.venv/bin/python bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stockpulse-bot

[Install]
WantedBy=multi-user.target
